# gromox-archive/supervisor.d/searchd.conf
[program:searchd]
command=/usr/bin/searchd --nodetach
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor-searchd.log
stderr_logfile=/var/log/supervisor-searchd-err.log
