# gromox-office/supervisor.d/crond.conf
[program:crond]
command=/usr/sbin/cron -n
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor-crond.log
stderr_logfile=/var/log/supervisor-crond-err.log
