# gromox-core/supervisor.d/gromox-http.conf
[program:gromox-http]
command=/usr/libexec/gromox/http
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor-gromox-http.log
stderr_logfile=/var/log/supervisor-gromox-http-err.log
