# gromox-core/supervisor.d/gromox-timer.conf
[program:gromox-timer]
command=/usr/libexec/gromox/timer
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor-gromox-timer.log
stderr_logfile=/var/log/supervisor-gromox-timer-err.log
