# gromox-core/supervisor.d/gromox-pop3.conf
[program:gromox-pop3]
command=/usr/libexec/gromox/pop3
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor-gromox-pop3.log
stderr_logfile=/var/log/supervisor-gromox-pop3-err.log
