# gromox-core/supervisor.d/gromox-midb.conf
[program:gromox-midb]
command=/usr/libexec/gromox/midb
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor-gromox-midb.log
stderr_logfile=/var/log/supervisor-gromox-midb-err.log
