# gromox-core/supervisor.d/gromox-zcore.conf
[program:gromox-zcore]
command=/usr/libexec/gromox/zcore
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor-gromox-zcore.log
stderr_logfile=/var/log/supervisor-gromox-zcore-err.log
