# gromox-core/supervisor.d/grommunio-admin-api.conf
[program:grommunio-admin-api]
command=/usr/sbin/uwsgi --ini /usr/share/grommunio-admin-api/api-config.ini
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor-grommunio-admin-api.log
stderr_logfile=/var/log/supervisor-grommunio-admin-api-err.log
