# gromox-core/supervisor.d/gromox-imap.conf
[program:gromox-imap]
command=/usr/libexec/gromox/imap
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor-gromox-imap.log
stderr_logfile=/var/log/supervisor-gromox-imap-err.log
