# gromox-core/supervisor.d/gromox-event.conf
[program:gromox-event]
command=/usr/libexec/gromox/event
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor-gromox-event.log
stderr_logfile=/var/log/supervisor-gromox-event-err.log
