# gromox-core/supervisor.d/saslauthd.conf
[program:saslauthd]
command=/usr/sbin/saslauthd -d -a pam
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor-saslauthd.log
stderr_logfile=/var/log/supervisor-saslauthd-err.log
